# ansible/roles/maintenance/templates/maintenance-daily.timer.j2
[Unit]
Description=Daily Maintenance Timer
Requires=maintenance-daily.service

[Timer]
OnCalendar=*-*-* {{ maintenance_daily_schedule }}:00
Persistent=true

[Install]
WantedBy=timers.target
