# ansible/roles/maintenance/templates/maintenance-weekly.timer.j2
[Unit]
Description=Weekly Maintenance Timer
Requires=maintenance-weekly.service

[Timer]
OnCalendar={{ maintenance_weekly_schedule }}
Persistent=true

[Install]
WantedBy=timers.target
