# ansible/roles/maintenance/templates/maintenance-weekly.service.j2
[Unit]
Description=Weekly Maintenance - Full OS Updates with Reboot
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/ansible-playbook \
  {{ maintenance_workdir }}/ansible/playbooks/maintenance-weekly.yml \
  --inventory {{ maintenance_workdir }}/ansible/inventory/hosts.yml \
  --extra-vars "target_host={{ ansible_hostname }}"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=maintenance-weekly
